# config_renderer/templates/units/k8s-addons.service
[Unit]
Description=Kubernetes Addons

[Service]
Type=oneshot
ExecStart=/opt/k8s-addons

[Install]
WantedBy=multi-user.target
