# src/seedbox/config_renderer/ignition/dnsmasq/dnsmasq.service
[Unit]
Description=dnsmasq DNS service

[Service]
Restart=always
RestartSec=10
ExecStart=/usr/bin/docker run \
    -p 53:53/tcp \
    -p 53:53/udp \
    --cap-add=NET_ADMIN \
    andyshinn/dnsmasq \
    --log-facility=- \
    {% for n in cluster.nodes.all() -%}
    --host-record={{ n.fqdn }},{{ n.ip }} \
    {% endfor %}

[Install]
WantedBy=multi-user.target
