# src/seedbox/config_renderer/ignition/etcd_server/etcd-member.service.d/30-image.conf
[Service]
Environment="ETCD_IMAGE_URL={{ config.etcd_image }}"
Environment="ETCD_IMAGE_TAG={{ cluster.etcd_image_tag }}"
