# src/seedbox/config_renderer/ignition/etcd_server/etcd-member.service.d/30-version.conf
[Service]
Environment="ETCD_IMAGE_TAG={{ config.etcd3_image_tag }}"
