# src/seedbox/config_renderer/ignition/etcd_server/etcd-member.service.d/30-version.conf
[Service]
Environment="ETCD_IMAGE_TAG={{ cluster.etcd_image_tag }}"
