# src/seedbox/config_renderer/ignition/etcd_server/etcd-member.service.d/40-persistent.conf
[Service]
Environment="ETCD_DATA_DIR={{ node.persistent_mountpoint.where }}/etcd"
