# src/seedbox/config_renderer/ignition/flannel/flanneld.service.d/40-network-config.conf
[Service]
ExecStartPre=/opt/bin/cluster-etcdctl --timeout 30s \
  set /coreos.com/network/config '{"Network": "{{ cluster.k8s_pod_network }}", "Backend": {"Type": "vxlan"}}'
