# src/seedbox/config_renderer/ignition/cni/docker.service.d/40-flannel.conf
[Unit]
Requires=flanneld.service
After=flanneld.service

[Service]
EnvironmentFile={{ config.k8s_cni_path }}/docker_opts_cni.env
