# src/seedbox/config_renderer/ignition/system/volume.mount
[Mount]
What={{ mountpoint.what }}
Where={{ mountpoint.where }}

[Install]
WantedBy={{ mountpoint.wanted_by }}
