# src/seedbox/config_renderer/ignition/rkt_runtime/rkt-api.service.d/30-proxy.conf
[Service]
Environment=HTTP_PROXY={{ cluster.aci_proxy_url }}
