# src/seedbox/config_renderer/ignition/etcd_server/locksmithd.service.d/40-etcd-lock.conf
[Service]
Environment="REBOOT_STRATEGY=etcd-lock"
